# packaging/com.zasupport.shield.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zasupport.shield</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/za-support-diagnostics/agent/za_shield_agent.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/zasupport-shield.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/zasupport-shield-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ZA_INSTALL_DIR</key>
        <string>/usr/local/za-support-diagnostics</string>
    </dict>
</dict>
</plist>
